# runecast.cron
#
# Regular cron jobs for RunecastUpdate
#
30 4    * * *   root            /usr/local/bin/RunecastUpdate
